# tb-localfolders-prune.service
[Unit]
Description=Thunderbird Local Folders restic prune

[Service]
Type=oneshot
ExecStart=%h/bin/thunderbird-localfolders-backup prune
