# tb-localfolders-backup.service
[Unit]
Description=Thunderbird Local Folders restic backup

[Service]
Type=oneshot
ExecStart=%h/bin/thunderbird-localfolders-backup
